# python/plugins/broker/mqtt/tests/docker_mqtt/container/broker/setup-pw.conf
[program:setup-pw]
command=/usr/local/bin/setup-pw.sh
autostart=true
autorestart=false
startsecs=0
startretries=0
stdout_logfile=/var/log/%(program_name)s.log
redirect_stderr=true
