# docs/tutorial/all-in-one/supervisord.conf
[unix_http_server]
file=/var/run/supervisor/supervisor.sock

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
nodaemon=true
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor/supervisor.sock

[program:setup-sshd]
command=bash -c "[ -f /etc/ssh/ssh_host_rsa_key ] || /usr/sbin/sshd-keygen"
autostart=true
startretries=1
priority=0
startsecs=0

[program:sshd]
command=/usr/sbin/sshd -D
autostart=%(ENV_ENABLE_SSHD)s
autorestart=true
stdout_logfile=/var/log/%(program_name)s.log
redirect_stderr=true

[program:zookeeper]
directory=/srv/kafka_2.12-%(ENV_KAFKA_VERSION)s
command=/srv/kafka_2.12-%(ENV_KAFKA_VERSION)s/bin/zookeeper-server-start.sh config/zookeeper.properties
autostart=%(ENV_ENABLE_BROKER)s
autorestart=true
stdout_logfile=/var/log/%(program_name)s.log
redirect_stderr=true

[program:kafka]
directory=/srv/kafka_2.12-%(ENV_KAFKA_VERSION)s
command=/srv/kafka_2.12-%(ENV_KAFKA_VERSION)s/bin/kafka-server-start.sh config/server.properties
autostart=%(ENV_ENABLE_BROKER)s
autorestart=true
stdout_logfile=/var/log/%(program_name)s.log
redirect_stderr=true

[program:mosquitto]
command=/usr/sbin/mosquitto -c /etc/mosquitto/mosquitto.conf
autostart=%(ENV_ENABLE_BROKER)s
autorestart=true
stdout_logfile=/var/log/%(program_name)s.log
redirect_stderr=true

