# python/plugins/broker/kafka/tests/docker_kafka/container/broker/setup-acl.conf

[program:setup-acl]
command=/usr/local/bin/setup-acl.sh
autostart=true
autorestart=false
startsecs=0
startretries=0
stdout_logfile=/var/log/%(program_name)s.log
redirect_stderr=true
