# java/plugin-kafka/src/test/docker/build2/setup-sasl.conf
[program:setup-sasl]
command=/usr/local/bin/setup-sasl.sh
autostart=true
autorestart=false
startsecs=0
startretries=0
stdout_logfile=/var/log/%(program_name)s.log
redirect_stderr=true

